# chrony.service
[Unit]
Description=Chrony NTP Container
Documentation=https://github.com/chris-short/chrony
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/home/cbshort/repo/chrony
User=cbshort
Group=cbshort

ExecStartPre=-/usr/bin/podman stop -t 10 chrony
ExecStartPre=-/usr/bin/podman rm chrony
ExecStartPre=/usr/bin/podman pull ghcr.io/chris-short/chrony:latest
ExecStart=/usr/bin/podman run --name=chrony \
    --network=host \
    --cap-add=SYS_TIME \
    -v chrony-data:/var/lib/chrony \
    -v chrony-logs:/var/log/chrony \
    -v /home/cbshort/repo/chrony/chrony.conf:/etc/chrony/chrony.conf:ro \
    ghcr.io/chris-short/chrony:latest
ExecStop=/usr/bin/podman stop -t 10 chrony
ExecStopPost=-/usr/bin/podman rm chrony

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
